# com.flighttracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.flighttracker</string>

  <key>ProgramArguments</key>
  <array>
    <!-- Update this path to match your Python 3.11 location -->
    <!-- Run: which python3.11   to find yours -->
    <string>/usr/local/bin/python3.11</string>
    <!-- Update this path to match where you cloned the repo -->
    <string>/Users/YOUR_USERNAME/flight-tracker-agent/main.py</string>
  </array>

  <key>WorkingDirectory</key>
  <!-- Update this to match where you cloned the repo -->
  <string>/Users/YOUR_USERNAME/flight-tracker-agent</string>

  <key>StartCalendarInterval</key>
  <array>
    <!-- 7:00 AM -->
    <dict>
      <key>Hour</key>
      <integer>7</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <!-- 7:00 PM -->
    <dict>
      <key>Hour</key>
      <integer>19</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </array>

  <key>StandardOutPath</key>
  <!-- Update this path too -->
  <string>/Users/YOUR_USERNAME/flight-tracker-agent/logs/tracker.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/YOUR_USERNAME/flight-tracker-agent/logs/tracker.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
